# dadbot.service
[Unit]
Description=Dadbot Discord bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dadbot
Group=dadbot
WorkingDirectory=/var/lib/dadbot

# Environment - now reads from file instead of exposing in argv
EnvironmentFile=/etc/dadbot.env
# No command line args needed - reads DISCORD_BOT_TOKEN from environment
ExecStart=/usr/local/bin/Dadbot

# Restart policy
Restart=on-failure
RestartSec=5s

# Logging (journald) with structured log identifiers
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dadbot
SyslogLevel=info
SyslogFacility=daemon

# Runtime directories (owned by User/Group automatically)
RuntimeDirectory=dadbot
StateDirectory=dadbot

# Security hardening (keep networking enabled for Discord API)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=

# Network security - only allow needed address families
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Hide process info from non-root (extra security)
ProtectProc=invisible

[Install]
WantedBy=multi-user.target
